# pr-review-bot.service
[Unit]
Description=PR Review Bot - AI-powered GitHub PR reviewer
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PR-REVIEW-XIBE
ExecStart=/usr/bin/node /home/ubuntu/PR-REVIEW-XIBE/bot.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pr-review-bot
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
